# .config/zsh/profile.d/datacube.sh
export DATACUBE_DB_URL='postgresql://mikael@localhost:5413/opendatacube'
export CFSI_OUTPUT_DIR="$HOME/files/cfsi_out"

