# .config/zsh/profile.d/datacube.sh
export DATACUBE_DB_URL='postgresql://opendatacube:opendatacube@localhost:5432/opendatacube'
export CFSI_BASE_DIR="$HOME/git/CFSI/cfsi"
export CFSI_CONTAINER_OUTPUT="$HOME/files/cfsi_container_out"
export CFSI_OUTPUT_DIR="$HOME/files/cfsi_out"

