# .config/zsh/profile.d/datacube.sh
export DATACUBE_DB_URL='postgresql://opendatacube:opendatacube@localhost:5432/opendatacube'
export CFSI_BASE_HOST="$HOME/git/CFSI/cfsi"
# export CFSI_BASE_CONTAINER="$HOME/git/CFSI/cfsi"
# export CFSI_CONFIG_FILE="config.yaml"
export CFSI_OUTPUT_TYPE=LOCAL
export CFSI_OUTPUT_HOST="$HOME/files/cfsi_out"
export CFSI_OUTPUT_CONTAINER="$HOME/files/cfsi_container_out"
