# .config/zsh/profile.d/editor.sh
export EDITOR=nvim
export VISUAL="idea -e"
