# .config/zsh/profile.d/mml-qgis-pistemeta-dev.sh
export PGDATABASE="pistedev01"
export PGPORT=5436
export PGHOST="localhost"
export PGSSLMODE="disable"
export QGIS_PLUGIN_USE_DEBUGGER="pydevd"
export PISTEMETA_LAZ_DIR="/home/mikael/Documents/gispo/pistemeta/44test/test_data/laz/automaattinen"

